# launchd/io.vectorize.hindsight.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.vectorize.hindsight.hourly</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.hindsight/venv/bin/python3</string>
        <string>__HOME__/.hindsight/nightly-learn.py</string>
        <string>--mode</string>
        <string>hourly</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.hindsight/logs/hourly-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.hindsight/logs/hourly-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.hindsight/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>GOOGLE_APPLICATION_CREDENTIALS</key>
        <string>__HOME__/.config/gcloud/application_default_credentials.json</string>
    </dict>
</dict>
</plist>
